# taskforge/deploy/supervisor.conf
; TaskForge Supervisor Configuration
; Place at: /etc/supervisor/conf.d/taskforge.conf
; Then: sudo supervisorctl reread && sudo supervisorctl update

[program:taskforge-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/taskforge/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/taskforge/storage/logs/worker.log
stopwaitsecs=3600
